# pwm-fan.service
[Unit]
Description=PWM fan control

[Service]
Type=simple
ExecStart=/usr/bin/start-pwm-fan /etc/pwm-fan.conf

[Install]
WantedBy=multi-user.target